# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
GenericName[ru]=Настройка цветовой температуры
GenericName[cs]=Přizpůsobení barevné teploty
Name[zh_CN]=红移
Name[es]=Redshift
Name[fr]=Redshift
GenericName[it]=Regolazione della temperatura del colore
Version=1.0
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Comment[da]=Justeringsværktøj for farvetemperatur
GenericName[fr]=Réglage de la température de couleur
Hidden=false
Comment[de]=Farbtemperatur Anpassungswerkzeug
Type=Application
X-GNOME-Autostart-enabled=true
Name=Redshift
GenericName[da]=Justering af farvetemperatur
StartupNotify=true
Exec=redshift-gtk
Comment[es]=Herramienta para el ajuste de la temperatura de color
GenericName[de]=Farbtemperatur Anpassung
GenericName[pl]=Dostosowanie temperatury barwowej
Categories=Utility;
Icon=redshift
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
Comment[it]=Strumento per la regolazione della temperatura del colore
GenericName[zh_CN]=色温调节
GenericName[es]=Ajuste de la temperatura de color
Name[tr]=Redshift
GenericName[tr]=Renk sıcaklığı ayarı
GenericName=Color temperature adjustment
Comment[tr]=Renk sıcaklığı ayarlama aracı
Name[it]=Redshift
GenericName[hu]=Szín hőmérséklet beállítás
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
Name[ru]=Redshift
Name[pl]=Redshift
Name[lt]=Redshift
Comment=Color temperature adjustment tool
Terminal=false
Name[pt_BR]=Redshift
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
Comment[zh_CN]=色温调节工具
GenericName[nl]=Kleurtemperatuur aanpassing
Comment[hu]=Színhőmérséklet beállító eszköz
Comment[fr]=Outil de réglage de la température de couleur
Name[de]=Redshift
Comment[ru]=Инструмент регулирования цветовой температуры
Name[hu]=Redshift
GenericName[lt]=Spalvos temperatūros reguliavimas

